# litespeed/etc/service/lsws/run
#!/bin/bash

exec 2>&1
source /env

exec /usr/local/lsws/bin/lswsctrl start
