# docker/tags/lighttpd/etc/service/lighttpd/run
#!/bin/sh

exec 2>&1
source /env

exec /usr/sbin/lighttpd -D -f /etc/lighttpd/example.conf
