# docker/tags/litespeed/etc/service/lsws/run
#!/bin/sh

exec 2>&1
source /env

exec /var/lib/openlitespeed/bin/lswsctrl start
