# litespeed/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:lsws]
command=/usr/local/lsws/bin/lswsctrl start
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
